# scripts/backoffice-api.service
[Unit]
Description=Back Office dashboard + API (loopback only)
Documentation=https://github.com/CodyJo/back-office/blob/main/docs/agents.md
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/projects/back-office
# `serve` is the full surface — dashboard static files + every /api/*
# route including the Phase 9 agent endpoints. Binds 127.0.0.1; the
# operator points admin.codyjo.com at this via localStorage["bo.api_base"].
ExecStart=/usr/bin/python3 -m backoffice serve --port 8071
# Restart on crash, but do not loop too fast.
Restart=on-failure
RestartSec=5
# Keep stdout/stderr in journald so `journalctl --user -u backoffice-api` works.
StandardOutput=journal
StandardError=journal
# Minimal env. The api-server reads config from BACK_OFFICE_CONFIG (or
# the default config/backoffice.yaml in WorkingDirectory).
Environment=PYTHONUNBUFFERED=1
Environment=BACK_OFFICE_ROOT=%h/projects/back-office

# Tighten the process slightly. Loopback bind + operator API key is the
# auth model; these are belt-and-suspenders.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=%h/projects/back-office/results %h/projects/back-office/dashboard %h/projects/back-office/config
ProtectHome=read-only
PrivateDevices=true
LockPersonality=true

[Install]
WantedBy=default.target
